# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[program:fake-dynamics]
command=dotnet /app/fake-dynamics/OrgReg.FakeDynamics.dll
environment=ASPNETCORE_URLS="http://+:5020",ASPNETCORE_ENVIRONMENT="Development"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10

[program:organisation-api]
command=dotnet /app/organisation-api/OrgReg.Organisation.Api.dll
environment=ASPNETCORE_URLS="http://+:5010",ASPNETCORE_ENVIRONMENT="Development",Dynamics__BaseUrl="http://localhost:5020",Dynamics__TenantId="fake",Dynamics__ClientId="fake",Dynamics__ClientSecret="fake",Dynamics__Resource="http://localhost:5020",Dynamics__ApiVersion="v9.2",Cors__Origins__0="*"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=20
startsecs=3

[program:frontend]
command=node /app/frontend/server.js
environment=PORT="3000",NODE_ENV="production",HOSTNAME="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=30
startsecs=3

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=40
startsecs=5
